# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:sonarqube]
command=/opt/startSonar.sh

[program:mysql]
command=mysqld_safe