# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:sonarqube]
command=sh /opt/sonarqube-4.5/bin/linux-x86-64/sonar.sh console